# com.obn.productivity-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.obn.productivity-tracker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>OBN_PI_INSTALL_DIR/obn_pi.py</string>
        <string>OBN_PI_VAULT_DIR</string>
        <string>--log-only</string>
    </array>
    <key>WorkingDirectory</key>
    <string>OBN_PI_INSTALL_DIR</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>55</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>~/Library/Logs/obn_pi.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/obn_pi.err</string>
</dict>
</plist>
